# AppDir/usr/lib/gsm/systemd/gsm.timer
[Unit]
Description=Run GSM automatic backup on schedule

[Timer]
OnCalendar=*-*-* 07..12:00
Persistent=true

[Install]
WantedBy=timers.target
