# AppDir/usr/lib/gsm/systemd/gsm.service
[Unit]
Description=GSM automatic backup service

[Service]
Type=oneshot
WorkingDirectory=%h/gsm
ExecStart=%h/gsm/src/gsm_cli.sh backup
